# whyisee-community/app/deploy/whyisee-community.service
[Unit]
Description=whyisee.xyz community app
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/whyisee-community/app
Environment=NODE_ENV=production
Environment=SITE_URL=https://whyisee.xyz
Environment=DB_HOST=127.0.0.1
Environment=DB_PORT=15432
Environment=DB_NAME=zi
Environment=DB_USER=postgres
Environment=DB_PASSWORD=123456
Environment=DB_SCHEMA=ws
Environment=PG_DUMP_BIN=pg_dump
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
